# deploy/launchd/com.ibshi.vattu.auditlogcleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- S2-4 Weekly Sunday 03:00 AuditLog cleanup (>90 days) -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ibshi.vattu.auditlogcleanup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/trinhhuuhung/Desktop/HUNGAI/HUNGTH OBSIDIAN V/HUNGTH OBSIDIAN/VẬT TƯ/scripts/audit_log_cleanup.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>3</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/vattu-audit-cleanup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/vattu-audit-cleanup.err</string>
</dict>
</plist>
